# build/supervisor-backend.ini
[program:backend]
directory=/src/
command=make backend
redirect_stderr=true
autostart=true
autorestart=true
startretries=3
stopwaitsecs=5
process_name=%(program_name)s
numprocs=1
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
